# config/services/aletheia.service
[Unit]
Description=Aletheia AI Gateway
After=network-online.target
Wants=network-online.target aletheia-memory.service

[Service]
Type=simple
EnvironmentFile=-%h/.aletheia/env
ExecStart=%h/.aletheia/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

# Usage:
#   aletheia-setup         # recommended — one command installs + starts everything
#
#   Manual install:
#   mkdir -p ~/.config/systemd/user
#   cp config/services/aletheia{,-memory}.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now aletheia-memory aletheia
#   loginctl enable-linger
#
# macOS launchd equivalent (~/.config/aletheia/com.aletheia.plist):
#
# <?xml version="1.0" encoding="UTF-8"?>
# <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
#   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
# <plist version="1.0">
# <dict>
#   <key>Label</key>
#   <string>com.aletheia</string>
#   <key>ProgramArguments</key>
#   <array>
#     <string>/bin/bash</string>
#     <string>/Users/<user>/.aletheia/start.sh</string>
#   </array>
#   <key>RunAtLoad</key>
#   <true/>
#   <key>KeepAlive</key>
#   <true/>
#   <key>StandardOutPath</key>
#   <string>/tmp/aletheia.log</string>
#   <key>StandardErrorPath</key>
#   <string>/tmp/aletheia.log</string>
# </dict>
# </plist>
#
# launchctl load ~/Library/LaunchAgents/com.aletheia.plist
