# config/services/aletheia.service
[Unit]
Description=Aletheia AI Gateway
After=network-online.target
Wants=network-online.target aletheia-memory.service

[Service]
Type=simple
Environment=ALETHEIA_HOME=__ALETHEIA_HOME__
EnvironmentFile=-__ALETHEIA_HOME__/instance/config/env
ExecStart=__NODE_BIN__ __ALETHEIA_HOME__/infrastructure/runtime/dist/entry.mjs gateway start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

# Usage:
#   aletheia enable            # recommended — registers for boot persistence
#   aletheia start             # start services now
#
#   Manual install:
#   mkdir -p ~/.config/systemd/user
#   cp config/services/aletheia{,-memory}.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now aletheia-memory aletheia
#   loginctl enable-linger
