# instance.example/services/aletheia.service
[Unit]
Description=Aletheia AI Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=RUST_BACKTRACE=1
# Optional: set API key here if not stored in instance/config/credentials/anthropic.json
# Environment=ANTHROPIC_API_KEY=sk-ant-...
# Load additional env vars from a file (silently ignored if the file does not exist).
# OAuth tokens written by the deploy script land here automatically.
EnvironmentFile=-%h/aletheia/instance/config/env
# %h expands to the user's home directory at runtime.
ExecStart=%h/.local/bin/aletheia -r %h/aletheia/instance
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

# --- Quick install ---
#
# 1. Copy this file:
#    mkdir -p ~/.config/systemd/user
#    cp instance.example/services/aletheia.service ~/.config/systemd/user/aletheia.service
#
# 2. If your binary is not at ~/.local/bin/aletheia, edit ExecStart, e.g.:
#    sed -i 's|%h/.local/bin/aletheia|/usr/local/bin/aletheia|g' ~/.config/systemd/user/aletheia.service
#
# 3. If your instance is not at ~/aletheia/instance, edit the ExecStart path, e.g.:
#    sed -i 's|%h/aletheia/instance|/srv/aletheia/instance|g' ~/.config/systemd/user/aletheia.service
#
# 4. Enable and start:
#    systemctl --user daemon-reload
#    systemctl --user enable --now aletheia
#
# 5. Persist across logout (run once per user):
#    loginctl enable-linger
#
# 6. Check logs:
#    journalctl --user -u aletheia -f
#
# 7. Verify:
#    aletheia health
