# instance.example/services/aletheia-health.timer
[Unit]
Description=Aletheia prosoche heartbeat timer
Requires=aletheia.service

[Timer]
OnActiveSec=60
OnUnitActiveSec=300
AccuracySec=30

[Install]
WantedBy=timers.target
