# instance.example/services/aletheia-health.service
[Unit]
Description=Aletheia health check
After=aletheia.service

[Service]
Type=oneshot
ExecStart=%h/aletheia/scripts/health-monitor.sh
StandardOutput=append:%h/aletheia/instance/logs/health.log
StandardError=append:%h/aletheia/instance/logs/health.log
