# instance.example/services/aletheia-health.service
[Unit]
Description=Aletheia prosoche heartbeat
Requires=aletheia.service
After=aletheia.service

[Service]
Type=oneshot
Environment=ALETHEIA_URL=http://127.0.0.1:18789
ExecStart=%h/.local/bin/aletheia-heartbeat
StandardOutput=journal
StandardError=journal
