# instance.example/services/com.aletheia.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aletheia.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>__ALETHEIA_ROOT__/target/release/aletheia</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>ALETHEIA_ROOT</key>
        <string>__ALETHEIA_ROOT__/instance</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ALETHEIA_ROOT__/instance/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>__ALETHEIA_ROOT__/instance/logs/gateway.err</string>
</dict>
</plist>
