# infrastructure/prosoche/aletheia-prosoche.service
[Unit]
Description=Aletheia Prosoche — Adaptive Attention Engine
After=network.target aletheia.service
Wants=aletheia.service

[Service]
Type=simple
User=syn
WorkingDirectory=<your-aletheia-root>/infrastructure/prosoche
ExecStart=<your-aletheia-root>/infrastructure/prosoche/.venv/bin/python3 -m prosoche.daemon <your-aletheia-root>/infrastructure/prosoche/config.yaml
Restart=on-failure
RestartSec=15
Environment=PATH=<your-aletheia-root>/shared/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
